# script/crontab
* * * * * bash /home/lucasmedeiros/Git/tcc/script/cpu_temperature.sh